# scripts/launchd/com.aionis.memory.api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.aionis.memory.api</string>

    <key>WorkingDirectory</key>
    <string>__AIONIS_ROOT__</string>

    <key>ProgramArguments</key>
    <array>
      <string>__AIONIS_ROOT__/scripts/daemon/run-api.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/aionis-memory-api.out.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/aionis-memory-api.err.log</string>
  </dict>
</plist>
